# etc/cron.daily/podcrto-backup
#!/bin/bash -e

chpst -u nobody:nogroup /code/backup.py
